# com.frozenturkey.locker.guard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
"http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.frozenturkey.locker.guard</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Library/Application Support/FrozenTurkeyLocker/guard.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Library/Application Support/FrozenTurkeyLocker/logs/guard-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Library/Application Support/FrozenTurkeyLocker/logs/guard-stderr.log</string>
</dict>
</plist>
